# rc.local

# Configure iptables to inspect bridged traffic
echo 1 >/proc/sys/net/bridge/bridge-nf-call-iptables
echo 1 >/proc/sys/net/bridge/bridge-nf-call-ebtables

ifconfig eth0 up 0.0.0.0
ifconfig eth1 up 0.0.0.0
brctl addbr br0
brctl addif br0 eth0 
brctl addif br0 eth1

udhcpc --now --interface br0

source /etc/firewall.user

exit 0
